# src/aws/bot.service
[Unit]
Description=WebullTradeAI trading bot
After=network-online.target docker.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=oneshot
User=root

# .env is mounted read-only so python-dotenv can parse it (handles quoted values
# correctly, unlike Docker's --env-file which passes them literally)
ExecStart=/usr/bin/docker run --rm \
    -v /home/ec2-user/WebullTradeAI/.env:/app/.env:ro \
    -v /home/ec2-user/WebullTradeAI/src/data:/app/src/data \
    -v /home/ec2-user/WebullTradeAI/src/models:/app/src/models:ro \
    webull-bot:latest

Restart=on-failure
RestartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webull-bot

# On weekdays, write Docker exit result to S3 (for crash diagnostics), then shut down.
# On weekends (Sat=6, Sun=7), retrain may be running — skip shutdown.
# Uses EC2 IAM role for S3 access — no .env credentials needed.
# $SERVICE_RESULT and $EXIT_CODE are substituted by systemd before execution.
ExecStopPost=/bin/bash -c 'export AWS_DEFAULT_REGION=us-east-1; D=$(date +%%Y-%%m-%%d); echo "result=$SERVICE_RESULT exit=$EXIT_CODE at $(date -u)" | aws s3 cp - "s3://webull-trade-ai/diagnostics/$D-bot-exit.txt" 2>/dev/null || true; day=$(date +%%u); if [ "$day" -ge 6 ]; then exit 0; fi; sudo /sbin/shutdown -h now'

[Install]
WantedBy=multi-user.target
