# src/aws/bot.service
[Unit]
Description=WebullTradeAI trading bot
After=network-online.target docker.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=oneshot
User=root

# .env is mounted read-only so python-dotenv can parse it (handles quoted values
# correctly, unlike Docker's --env-file which passes them literally)
ExecStart=/usr/bin/docker run --rm \
    -v /home/ec2-user/WebullTradeAI/.env:/app/.env:ro \
    -v /home/ec2-user/WebullTradeAI/src/data:/app/src/data \
    -v /home/ec2-user/WebullTradeAI/src/models:/app/src/models:ro \
    webull-bot:latest

Restart=on-failure
RestartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webull-bot

# On weekdays, write Docker exit result to S3 (for crash diagnostics), then shut down.
# On weekends (Sat=6, Sun=7), retrain may be running — skip shutdown.
# $SERVICE_RESULT and $EXIT_CODE are set by systemd: result is success/exit-code/signal/timeout.
ExecStopPost=/bin/bash -c '\
    ENV=/home/ec2-user/WebullTradeAI/.env; \
    KEY=$(grep "^AWS_ACCESS_KEY_ID=" "$ENV" 2>/dev/null | cut -d= -f2- | tr -d "\""); \
    SEC=$(grep "^AWS_SECRET_ACCESS_KEY=" "$ENV" 2>/dev/null | cut -d= -f2- | tr -d "\""); \
    REG=$(grep "^AWS_REGION=" "$ENV" 2>/dev/null | cut -d= -f2- | tr -d "\""); \
    BKT=$(grep "^AWS_S3_BUCKET=" "$ENV" 2>/dev/null | cut -d= -f2- | tr -d "\""); \
    export AWS_ACCESS_KEY_ID="$KEY" AWS_SECRET_ACCESS_KEY="$SEC" AWS_DEFAULT_REGION="$REG"; \
    D=$(date +%%Y-%%m-%%d); \
    [ -n "$BKT" ] && [ -n "$KEY" ] && \
        echo "result=$SERVICE_RESULT exit=$EXIT_CODE at $(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ)" | \
        aws s3 cp - "s3://$BKT/diagnostics/$D-bot-exit.txt" 2>/dev/null || true; \
    day=$(date +%%u); if [ "$day" -ge 6 ]; then exit 0; fi; sudo /sbin/shutdown -h now'

[Install]
WantedBy=multi-user.target
