# src/aws/bot.timer
[Unit]
Description=Run WebullTradeAI bot daily at midday
After=config-pull.service

[Timer]
# 11:00 AM ET = 8:00 AM PT = 15:00 UTC
OnCalendar=Mon..Fri 15:00:00 UTC
AccuracySec=60

[Install]
WantedBy=timers.target
