# src/aws/retrain.service
[Unit]
Description=WebullTradeAI weekly model retrain
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root

ExecStart=/usr/bin/docker run --rm \
    --gpus all \
    -v /home/ec2-user/WebullTradeAI/.env:/app/.env:ro \
    -v /home/ec2-user/WebullTradeAI/src/data:/app/src/data \
    -v /home/ec2-user/WebullTradeAI/src/models:/app/src/models \
    webull-retrain:latest

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webull-retrain

# Kill the container if retrain exceeds 8 hours — prevents runaway charges.
TimeoutStartSec=8h
TimeoutStopSec=60

# Always shut down after retrain (success, failure, or timeout).
# Single command with no bash logic — nothing that can fail or syntax-error.
ExecStopPost=/sbin/shutdown -h now

[Install]
WantedBy=multi-user.target
