# src/aws/retrain.service
[Unit]
Description=WebullTradeAI weekly model retrain
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root

ExecStart=/usr/bin/docker run --rm \
    --gpus all \
    -v /home/ec2-user/WebullTradeAI/.env:/app/.env:ro \
    webull-retrain:latest

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webull-retrain

# Retrain can take several hours — no timeout
TimeoutStartSec=infinity

# Shutdown is handled by retrain.py itself on success to avoid idle EC2 costs.

[Install]
WantedBy=multi-user.target
