# src/aws/retrain.service
[Unit]
Description=WebullTradeAI weekly model retrain
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root

ExecStart=/usr/bin/docker run --rm \
    --gpus all \
    -v /home/ec2-user/WebullTradeAI/.env:/app/.env:ro \
    -v /home/ec2-user/WebullTradeAI/src/data:/app/src/data \
    -v /home/ec2-user/WebullTradeAI/src/models:/app/src/models \
    webull-retrain:latest

StandardOutput=journal
StandardError=journal
SyslogIdentifier=webull-retrain

# Kill the container and trigger ExecStopPost (EC2 shutdown) if retrain
# exceeds 8 hours — 2× the expected ~4h runtime. Prevents runaway charges.
TimeoutStartSec=8h
TimeoutStopSec=60

# Shut down after retrain completes (success or failure) on weekends.
# retrain.py ran shutdown internally pre-Docker, but that doesn't work inside
# a container — ExecStopPost runs on the EC2 host after the container exits.
ExecStopPost=/bin/bash -c 'day=$(date +%%u); if [ "$day" -lt 6 ]; then exit 0; fi; sudo /sbin/shutdown -h now'

[Install]
WantedBy=multi-user.target
